# .config/autostart/cmst.desktop
[Desktop Entry]
Exec=cmst -m
Name=cmst
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
