# .config/autostart/nm-applet.desktop
[Desktop Entry]
Comment=Networkmmanager applet
Exec=nm-applet --indicator
Hidden=true
Icon=network-wireless
Name=Network
Terminal=false
Type=Application
X-LXQt-Module=true
