# .config/autostart/qps.desktop
[Desktop Entry]
Exec=qps -mini
Name=qps
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
